# runit/services/x11vnc/run
#!/bin/sh

exec x11vnc -noxdamage -shared -many -auth ~/.Xauthority -display $DISPLAY
