# systemd/system/checkupdates.timer
[Unit]
Description=Check if updates system updates are available

[Timer]
OnActiveSec=5m
Unit=checkupdates.service

[Install]
WantedBy=timers.target
