# runit/services/picom/run
#!/bin/sh

exec picom --experimental-backends
