# runit/services/polybar/run
#!/bin/sh

killall -9 polybar

polybar --reload mainbar-i3 &
exec polybar --reload altbar-i3 
