# runit/services/polybar/run
#!/bin/sh

killall -9 polybar

exec polybar --reload mainbar-i3
