# runit/services/xclip/run
#!/bin/sh

xclip -i /tmp/gscreenshot-cli-clip.png -selection clipboard -t image/png -quiet
