# runit/services/dbus/run
#!/bin/sh

exec dbus-daemon --nofork --nopidfile --session
