# runit/services/pulseaudio/run
#!/bin/sh

/usr/bin/pulseaudio --daemonize=no
